# letsencrypt/letsencrypt.service.d/haproxy.conf
[Service]
ExecStartPost=/usr/bin/haproxy-letsencrypt.sh
