# letsencrypt/letsencrypt.timer
[Unit]
Description=letsencrypt timer

[Timer]
OnCalendar=daily
Persistent=true
Unit=letsencrypt.service

[Install]
WantedBy=basic.target
